# update-git-repository@.timer
[Unit]
Description=Update a Git repository daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
